# niwalaps.timer
[Unit]
Description=Run niwalaps every morning

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
